# content-update-bot.service
[Unit]
Description=Content Update Bot
After=network.target

[Service]
Type=simple
User=vincent
Group=vincent
WorkingDirectory=/home/vincent/Documents/Code/BotBlogScraper
Environment=PYTHONPATH=/home/vincent/Documents/Code/BotBlogScraper
ExecStart=/home/vincent/Documents/Code/BotBlogScraper/venv/bin/python main.py
Restart=always
RestartSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ProtectHome=read-only
ReadWritePaths=/home/vincent/Documents/Code/BotBlogScraper/storage /home/vincent/Documents/Code/BotBlogScraper/bot.log

[Install]
WantedBy=multi-user.target